# example/etc/cron.daily/fsvs-versioning
#!/bin/sh

/var/lib/fsvs-versioning/scripts/commit.sh "Commit per cron.daily"
